# releases/v4.3/files/etc/systemd/system/repeater-time-announce.service
[Unit]
Description=Prepare hourly time announcement for repeater

[Service]
Type=oneshot
ExecStart=/usr/local/bin/repeater_time_announce.py
