# releases/v4.7/files/etc/systemd/system/repeater-selfheal.service
[Unit]
Description=ProtoRadio autocura de rede e servicos
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/repeater-selfheal

